# siftengine.service
[Unit]
Description=SiftEngine log processing pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/Users/sauravbhatta/Desktop/SiftEngine
ExecStart=/Users/sauravbhatta/Desktop/SiftEngine/sift.sh
EnvironmentFile=-/etc/siftengine/siftengine.env

[Install]
WantedBy=multi-user.target
